# core/unittest/cms/unittest_data/other/conf/cloudmonitor.usr
#!/bin/bash
#chkconfig: 2345 80 90
bash BASE_PATH/cloudmonitorCtl.sh $1-directly